# roles/common/templates/docker-etcd/supervisord.conf
[program:etcd]
command=/usr/local/bin/run_etcd.sh
