# roles/postgres/templates/docker-postgres/supervisord.conf
[program:postgres]
user=postgres
group=postgres
command=/usr/bin/postmaster -D "{{ postgres.data_dir }}"
autostart=true
autorestart=true
redirect_stderr=true
stopsignal=QUIT